# sensors/sensors.service
[Unit]
Description=Sensors

[Service]
ExecStart=/bin/bash -c "cd /opt/sensors; ./venv/bin/python app.py"

[Install]
WantedBy=multi-user.target
